# .config/systemd/user/bluetooth-reset.service
[Unit]
Description=Bluetooth reset
After=graphical-session.target
Requires=graphical-session.target

[Service]
Type=oneshot
ExecStart=/usr/bin/sh -c 'sudo systemctl restart bluetooth && sudo rfkill block bluetooth && sudo rfkill unblock bluetooth'
RestartSec=0

[Install]
WantedBy=graphical-session.target
