# .config/systemd/user/bluetooth-reset.service
[Unit]
Description=Bluetooth reset
After=graphical-session.target
Requires=graphical-session.target

[Service]
Type=oneshot
ExecStart=/usr/bin/sudo modprobe -r btusb
ExecStart=/usr/bin/sudo modprobe btusb
ExecStart=/usr/bin/sleep 1
ExecStart=/usr/bin/bluetoothctl power on

[Install]
WantedBy=graphical-session.target
