# .config/systemd/user/daily.timer
[Unit]
Description=Start Daily
PartOf=work.target

[Timer]
OnCalendar=Mon..Fri, 09:30
Persistent=true
Unit=google-chrome@daily.service
