# .config/systemd/user/deluge-arch-torrent.timer
[Unit]
Description=Arch torrent timer

[Timer]
OnCalendar=*-*-02
AccuracySec=1m
Persistent=true

[Install]
WantedBy=graphical-session.target
