# .config/systemd/user/filebrowser-proxy@.socket
[Socket]
ListenStream=127.0.0.1:$PORT

[Install]
WantedBy=sockets.target
