# .config/systemd/user/gamemode.target
[Unit]
Description=gamemode
Requires=gamemode.service
Wants=gamemode.service
Conflicts=work.target
