# .config/systemd/user/glava.target
[Unit]
Description=glava bundle
Requires=graphical-session.target pulseaudio.socket
After=graphical-session.target pulseaudio.socket
Wants=glava@bar.service glava@radial.service glava-update.service

[Install]
WantedBy=graphical-session.target
