# .config/systemd/user/glava.target
[Unit]
Description=glava
Requires=graphical-session.target pipewire-pulse.socket
After=graphical-session.target pipewire-pulse.socket
Wants=glava@bar.service glava-update.service

[Install]
WantedBy=graphical-session.target
