# .config/systemd/user/glava.target
[Unit]
Description=glava
Requires=i3-session.target pipewire-pulse.socket
After=i3-session.target pipewire-pulse.socket
Wants=glava@bars.service glava-update.service

[Install]
WantedBy=i3-session.target
