# .config/systemd/user/i3-session.target
[Unit]
BindsTo=xorg-session.target
After=xorg-session.target
PropagatesStopTo=xorg-session.target
BindsTo=i3.service
After=i3.service
