# .config/systemd/user/i3-session.target
[Unit]
Description=i3
BindsTo=graphical-session.target
Wants=graphical-session-pre.target
After=graphical-session-pre.target
