# .config/systemd/user/i3.service
[Unit]
Requires=xorg.service
After=xorg.service
BindsTo=i3-session.target

[Service]
ExecStart=/usr/bin/i3
Nice=-15
