# .config/systemd/user/jaeger-otlp-proxy.socket
[Socket]
ListenStream=127.0.0.1:4317
Service=localhost-proxy@4317.service
