# .config/systemd/user/joystickwake.service
[Unit]
StopWhenUnneeded=true
RefuseManualStart=true
RefuseManualStop=true

[Service]
Type=exec
ExecStart=/usr/bin/joystickwake
Restart=always
