# .config/systemd/user/kanshi.service
[Unit]
Description=kanshi
After=sway-session.target
Requires=sway-session.target

[Service]
Type=exec
ExecStart=/usr/bin/kanshi
Restart=always
