# .config/systemd/user/kubernetes-cluster-proxy.socket
[Socket]
ListenStream=127.0.0.1:6443

[Install]
WantedBy=sockets.target
