# .config/systemd/user/localhost-proxy@.service
[Service]
Environment=TIMEOUT=10m
ExecStart=/usr/lib/systemd/systemd-socket-proxyd --exit-idle-time=${TIMEOUT} 127.0.0.2:%I
