# .config/systemd/user/mariadb-proxy.socket
[Socket]
ListenStream=127.0.0.1:3306
Service=localhost-proxy@3306.service
