# .config/systemd/user/messengers.target
[Unit]
Description=Messengers
Requires=graphical-session.target
After=graphical-session.target
Wants=discord.service
Wants=google-chrome@whatsapp.service
Wants=telegram.service
