# .config/systemd/user/messengers.slice
[Unit]
Description=Messengers Slice

[Slice]
