# .config/systemd/user/messengers.target
[Unit]
Description=Messengers
Requires=graphical-session.target
After=graphical-session.target

[Install]
WantedBy=graphical-session.target
