# .config/systemd/user/messengers.target
[Unit]
Description=Messengers
Requires=graphical-session.target
After=graphical-session.target
Wants=discord.service
Wants=whatsapp.service
#Wants=telegram.service

[Install]
WantedBy=graphical-session.target
