# .config/systemd/user/podman-cleanup.service
[Unit]
Description=Clean podman storage

[Service]
Type=oneshot
ExecStart=/usr/bin/podman system prune --all --force --volumes --external
SuccessExitStatus=0 125
