# .config/systemd/user/podman-cleanup.timer
[Unit]
Description=Clean podman storage

[Timer]
OnCalendar=weekly
AccuracySec=1h
Persistent=true
