# .config/systemd/user/polybar@.service
[Unit]
Description=Polybar %I
After=graphical-session.target
Requires=graphical-session.target

[Service]
Type=forking
ExecStart=%h/.config/polybar/launch.sh %i
RestartSec=0
Restart=always

[Install]
WantedBy=graphical-session.target
