# .config/systemd/user/real-daily.timer
[Unit]
Description=Start Real Daily
PartOf=work.target

[Timer]
OnCalendar=Mon..Fri, 10:00
Persistent=true
Unit=real-daily-trigger.service
