# .config/systemd/user/restart-service@.service
[Unit]
Description=Restart service %I

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl --user restart '%i*.service'

[Install]
WantedBy=graphical-session.target
