# .config/systemd/user/rss.timer
[Unit]
Description=Start RSS Feeder
PartOf=work.target

[Timer]
OnCalendar=Thu 09:00
Persistent=true

[Install]
WantedBy=work.target
