# .config/systemd/user/sshuttle.slice
[Unit]
Description=sshuttle
Requisite=work.target

[Slice]