# .config/systemd/user/sshuttle.slice
[Unit]
Description=sshuttle Slice
Requisite=work.target

[Slice]