# .config/systemd/user/startup-boost.service
[Unit]
Requires=powerprofile@performance.service

[Service]
Type=exec
ExecStart=/bin/sleep 1m
RemainAfterExit=false
Restart=no
