# .config/systemd/user/steam.service
[Unit]
Description=Steam
After=graphical-session.target
Requires=graphical-session.target

[Service]
Type=simple
ExecStartPre=-/usr/bin/steam-tweaks
ExecStart=/usr/bin/steam
OOMScoreAdjust=-500
TimeoutStopSec=5s
MemoryMax=40%
TasksMax=800
RestartSec=0
Restart=always
ManagedOOMPreference=avoid
Nice=-2

[Install]
WantedBy=graphical-session.target
