# .config/systemd/user/steam.service
[Unit]
Description=Steam
After=graphical-session.target
Requires=graphical-session.target
StopWhenUnneeded=yes

[Service]
Type=simple
ExecStartPre=-/usr/bin/steam-tweaks
ExecStart=/usr/bin/steam
MemoryMax=80%
TasksMax=800
RestartSec=0
Restart=always

[Install]
WantedBy=graphical-session.target