# .config/systemd/user/stop-telefondienst.timer
[Unit]
Description=Stop Telefondienst
PartOf=work.target

[Timer]
OnCalendar=Thu 12:30
Persistent=true

[Install]
WantedBy=graphical-session.target
