# .config/systemd/user/stop-telefondienst.timer
[Unit]
Description=Stop Telefondienst
PartOf=telefondienst.target

[Timer]
OnCalendar=Thu 12:30
Persistent=true
