# .config/systemd/user/sway-session.target
[Unit]
Description=sway session
BindsTo=wayland-session.target
Wants=swaybg.service

