# .config/systemd/user/teams.service
[Unit]
Description=teams
PartOf=work.target

[Service]
Type=forking
ExecStart=/usr/bin/teams-insiders
RestartSec=0
Restart=always

