# .config/systemd/user/telefondienst.timer
[Unit]
Description=Start Telefondienst
PartOf=work.target

[Timer]
OnCalendar=Thu 09:00
Persistent=true
Unit=telefondienst.target
