# .config/systemd/user/telegram.service
[Unit]
Description=Telegram
PartOf=messengers.target
Requisite=graphical-session.target org.freedesktop.IBus.session.generic.service
After=graphical-session.target dunst.service
Wants=dunst.service

[Service]
Type=simple
Environment=QT_IM_MODULE=ibus
ExecStart=/usr/bin/telegram-desktop
RestartSec=0
Restart=always
Slice=messengers.slice

# disallow writing to /usr, /bin, /sbin, ...
ProtectSystem=true
# more paranoid security settings
NoNewPrivileges=true
ProtectKernelTunables=true
ProtectControlGroups=true
