# .config/systemd/user/unclutter.service
[Unit]
Description=Auto-hide mouse cursor
PartOf=i3-session.target

[Service]
Type=exec
ExecStart=/usr/bin/unclutter
