# .config/systemd/user/updates.timer
[Unit]
Description=Package Updates

[Timer]
OnCalendar=hourly
AccuracySec=1m
Persistent=true

[Install]
WantedBy=graphical-session.target
