# .config/systemd/user/updates.slice
[Unit]
Description=Updates Slice

[Slice]
CPUWeight=1
