# .config/systemd/user/updates.timer
[Unit]
Description=Updates

[Timer]
OnCalendar=hourly
AccuracySec=1m
Persistent=true
