# .config/systemd/user/work.target
[Unit]
Description=Work applications
Requires=work.target
After=graphical-session.target
Wants=work.target rocketchat.service google-chrome@work.service google-chrome@work-todo.service telefondienst.timer stop-telefondienst.timer daily.timer

[Install]
WantedBy=graphical-session.target
