# .config/systemd/user/work-stop.timer
[Unit]
Description=Stop Work
#PartOf=work.target # TODO: See https://github.com/systemd/systemd/pull/18517

[Timer]
OnActiveSec=8h
