# .config/systemd/user/work.target
[Unit]
Description=Work applications
Requires=work.target
After=graphical-session.target
Wants=work.target rocketchat.service google-chrome@work.service rss.timer stop-rss.timer

[Install]
WantedBy=graphical-session.target
