# .config/systemd/user/work.target
[Unit]
Description=Work applications
Requires=graphical-session.target
After=graphical-session.target
Wants=teams.service teams-audio-fix.service vpn@crossmedia.service

[Install]
WantedBy=graphical-session.target
