# .config/systemd/user/work.target
[Unit]
Description=Work applications
After=graphical-session.target
Wants=daily.timer
Wants=google-chrome@work-todo.service
Wants=google-chrome@work.service
Wants=real-daily.timer
Wants=rocketchat.service
Wants=stop-telefondienst.timer
Wants=telefondienst.timer
Wants=update-work-git-repos.service
Wants=work-stop-trigger.service
Wants=work.timer
#Wants=work-stop.timer # TODO: See https://github.com/systemd/systemd/pull/18517
