# .config/systemd/user/xorg-configs.timer
[Unit]
Description=Set xorg settings

[Timer]
OnActiveSec=10
AccuracySec=1s
Persistent=true

[Install]
WantedBy=xorg-session.target
