# .config/systemd/user/xorg-settings.timer
[Unit]
Description=Set xorg settings
After=i3-session.target

[Timer]
OnStartupSec=10
AccuracySec=1s
Persistent=true
