# .config/systemd/user/youtube-music-player.service
[Unit]
Description=YouTube Music Player
After=graphical-session.target pipewire.service
Requires=graphical-session.target pipewire.service
ConditionACPower=true

[Service]
Type=simple
ExecStart=/usr/bin/youtube-music-desktop-app
RestartSec=0
Restart=always
MemoryMax=1G
Nice=-10
ManagedOOMMemoryPressure=kill
ManagedOOMMemoryPressureLimit=10%

[Install]
WantedBy=graphical-session.target
