# .config/systemd/user/youtube-music-player.service
[Unit]
Description=YouTube Music Player
After=graphical-session.target network-online.service pipewire.service
Requires=graphical-session.target network-online.service pipewire.service

[Service]
Type=simple
ExecStart=/usr/bin/ytmdesktop
#ExecStart=/usr/bin/google-chrome-stable --user-data-dir=%E/youtube-music --app=https://music.youtube.com --class=youtube-music-desktop-app
RestartSec=0
Restart=always
MemoryMax=2G

[Install]
WantedBy=graphical-session.target
