# .config/systemd/user/youtube-music-player.service
[Unit]
Description=YouTube Music Player
Wants=pipewire-pulse.service
Requires=graphical-session.target
After=graphical-session.target pipewire-pulse.service
ConditionACPower=true

[Service]
Type=exec
ExecStart=/usr/bin/youtube-music-desktop-app
Restart=always
MemoryHigh=1500M
MemoryMax=2G
Nice=-10
ManagedOOMMemoryPressure=kill
ManagedOOMMemoryPressureLimit=10%
