# .config/systemd/user/youtube-music-player.service
[Unit]
Description=YouTube Music Player
After=graphical-session.target network-online.service
Requires=graphical-session.target network-online.service

[Service]
Type=simple
ExecStart=/usr/bin/ytmdesktop
RestartSec=0
Restart=always

[Install]
WantedBy=graphical-session.target
