# .config/systemd/user/localhost-proxy@8050.service.d/override.conf
[Unit]
Requires=filebrowser.service
After=filebrowser.service
