# config/rootfs/etc/systemd/system/systemd-fsck-@.service.d/override.conf
[Service]
StandardOutput=journal
StandardError=journal+console