# .config/systemd/user/google-chrome@stuff.service.d/override.conf
[Unit]
Description=Google Chrome
[Service]
Nice=-5
