# .config/systemd/user/gopass-fuse.service.d/override.conf
[Unit]
Requisite=work.target
RefuseManualStart=true
RefuseManualStop=true
