# .config/systemd/user/bottles-.service.d/override.conf
[Service]
ExecStop=pkill --cgroup %n -9 winedevice.exe
TimeoutStopSec=15
