# .config/systemd/user/filebrowser@files.service.d/override.conf
[Service]
Environment=ROOT=%h
Environment=PORT=8051