# config/rootfs/etc/systemd/system/user@.service.d/override.conf
[Service]
Delegate=yes
ManagedOOMMemoryPressure=kill
ManagedOOMMemoryPressureLimit=25%
