# config/rootfs/etc/systemd/system/avahi-daemon.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/avahi-daemon -s -f /etc/avahi/stupid-avahi-daemon.conf
