# config/rootfs/etc/udev/rules.d/50-xorg-settings.rules
SUBSYSTEM=="input", ACTION=="add", RUN+="/usr/local/bin/reload-xorg-settings"
