# .config/systemd/user/polybar@top.service.d/override.conf
[Unit]
Requires=pipewire-pulse.service
After=pipewire-pulse.service
