# .config/systemd/user/polybar@top.service.d/override.conf
[Unit]
Requires=pipewire.service
After=pipewire.service
