# .config/systemd/user/google-chrome@daily.service.d/override.conf
[Unit]
Requisite=work.target

[Service]
Environment=EXTRA_ARGS=--app=https://3cx.teuto.net/meet/4478b570b21285acd50a3f99378f4f29031ffaef
Restart=on-failure