# .config/systemd/user/google-chrome@daily.service.d/override.conf
[Unit]
Requisite=work.target
Requires=meeting-audio.service

[Service]
Environment=EXTRA_ARGS=--app=https://3cx.teuto.net/meet/82214e8a663bf1373fda9e8cd0133de281802573
Restart=on-failure
