# .config/systemd/user/google-chrome@daily.service.d/override.conf
[Unit]
Requisite=work.target

[Service]
Environment=EXTRA_ARGS=--app=https://3cx.teuto.net/meet/f9c5305fc005b8f26b5e6e2a03aa353947dc42ba
Restart=on-failure