# .config/systemd/user/google-chrome@daily.service.d/override.conf
[Unit]
Requisite=work.target

[Service]
Environment=EXTRA_ARGS=--app=https://3cx.teuto.net/meet/9d95d2e316dcce298003d664b4f90f0072aa182a
Restart=on-failure