# .config/systemd/user/google-chrome@daily.service.d/override.conf
[Unit]
Requisite=work.target

[Service]
Environment=EXTRA_ARGS=--app=https://3cx.teuto.net/meet/3bffe42ec914e254cd6387ae8efee9226feb7e2f
Restart=on-failure