# .config/systemd/user/clipmenud.service.d/override.conf
[Unit]
Requires=graphical-session.target

[Service]
RestartSec=10s
StartLimitBurst=50
Environment=CM_MAX_CLIPS=0
Environment=CM_DIR=%C/clipmenud