# .config/systemd/user/clipmenud.service.d/override.conf
[Unit]
Requires=graphical-session.target
