# .config/systemd/user/clipmenud.service.d/override.conf
[Service]
Environment=CM_MAX_CLIPS=0 CM_DIR=%h/.cache/clipmenud
