# config/rootfs/etc/systemd/system/paccache.service.d/override.conf
[Service]
ExecStart=/usr/bin/paccache -ruk1
