# .config/systemd/user/google-chrome@ytm.service.d/override.conf
[Unit]
Description=YouTube Music Player
Wants=pipewire.service
After=pipewire.service
ConditionACPower=true

[Service]
Environment=EXTRA_ARGS=--app=https://music.youtube.com
