# .config/systemd/user/browser@telefondienst.service.d/override.conf
[Unit]
PartOf=telefondienst.target

[Service]
Environment=EXTRA_ARGS=--app=https://gitlab.teuto.net/technik/teutoops/-/blob/master/docs/source/organisation/allgemein/telefondienst.md