# .config/systemd/user/redshift.service.d/override.conf
[Service]
ExecStartPre=/bin/bash -c 'until i3-msg; do :; done'
