# .config/systemd/user/localhost-proxy@4317.service.d/override.conf
[Unit]
Requires=jaeger.service
After=jaeger.service
[Service]
Environment=TIMEOUT=1m
