# .config/systemd/user/gamemoded.service.d/override.conf
[Unit]
Wants=gamemode.service
Requires=gamemode.service
