# config/rootfs/etc/systemd/system/reflector.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/reflector @/etc/xdg/reflector/reflector.conf @/etc/xdg/reflector/reflector.conf.d/override.conf
