# .config/systemd/user/filebrowser-proxy@files.service.d/override.conf
[Unit]
Wants=filebrowser@files.service
After=filebrowser@files.service

[Service]
Environment=PORT=8051