# .config/systemd/user/localhost-proxy@16686.service.d/override.conf
[Unit]
Requires=jaeger.service
After=jaeger.service
