# HOME/.config/systemd/user/xorg.socket
[Unit]
Description=Socket for xorg at display 0

[Socket]
ListenStream=/tmp/.X11-unix/X0
